# cznic/sfpswitch/files/sfpled
#!/bin/sh /etc/rc.common

START=99

start_service() {
  /usr/sbin/sfpswitch.py --resetled
}

