# hardware/userspace_time_sync/files/time_sync.cron
MAILTO=""
33	*	*	*	*	root	/sbin/hwclock -u -w -f /dev/rtc0
