# cznic/rainbow-omnia/files/rainbow.cron
MAILTO=""
*	*	*	*	*	root	/usr/bin/rainbow_button_sync.sh
