# cznic/turris-diagnostics/files/country-check.cron
## crontab

MAILTO=""
0 8 * * 1 root /usr/share/diagnostics/country-check.sh
