# cznic/user_notify/files/user_notify.cron
MAILTO=""
*/5     *       *       *       *       root    /usr/bin/notifier

