# hardware/mox/mox-shutdown-powersave/files/init
#!/bin/sh /etc/rc.common

STOP=99

start() {
	:
}

stop() {
	echo powersave > /sys/devices/system/cpu/cpufreq/policy0/scaling_governor
}

restart() {
	:
}

reload() {
	:
}
