# cznic/haas-mitmproxy/files/haas.cron
MAILTO=""

30	4	*	*	*	root	/usr/share/haas/register.sh
