# cznic/server-uplink/files/server-uplink.cron
MAILTO=""

0	*	*	*	*	root	/usr/share/server-uplink/registration_code.sh
