# cznic/server-uplink/files/server-uplink.cron
MAILTO=""

*/60	*	*	*	*	root	/usr/share/server-uplink/registration_code.sh
*/60	*	*	*	*	root	/usr/share/server-uplink/contract_valid.sh
