# cznic/turris-cagen/files/cagen.cron
MAILTO=""
1	1	*	*	*	root	/usr/bin/turris-cagen refresh
