# cznic/nethist/files/nethist.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2013,2016 CZ.NIC z.s.p.o. (http://www.nic.cz/)

START=50
USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command nethist
	procd_set_param respawn retry=60
	procd_close_instance
}
