# net/netmetr/files/netmetr.cron
MAILTO=""
00 	* 	* 	* 	* 	root	netmetr --autostart --dwlhist --syslog &
