# cznic/mox-support/files/init
#!/bin/sh /etc/rc.common

START=15
start() {
	/usr/share/mox/mox_autosetup
}

stop() {
	:
}
