# admin/vixie-cron/files/cron.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org
# Copyright (C) 2013-2016 CZ.NIC

START=90

USE_PROCD=1

start_service() {
	mkdir -p /var/spool/cron/
	if [ -e /var/spool/cron/crontabs ] ; then
		rm -rf /var/spool/cron/crontabs
	fi
	ln -sf /etc/crontabs /var/spool/cron
	procd_open_instance
	procd_set_param command /usr/sbin/cron -n
	procd_set_param respawn
	procd_close_instance
}
