# net/suricata/files/suricata.cron
15 * * * * root /usr/bin/suricata_update_rules.sh
