# cznic/foris-controller-ssbackups-module/files/ssbackups.cron
MAILTO=""

45	2	*	*	*	root	sleep $(tr -cd 0-9 </dev/urandom | head -c 2) && /usr/bin/foris-client -m ssbackups -a create_and_upload ubus
