# cznic/watchdog_adjust/files/etc/init.d/watchdog_adjust
#!/bin/sh /etc/rc.common

START=11

start() {
	ubus call system watchdog '{"magicclose": true}'
	ubus call system watchdog '{"stop": true}'
	watchdog -t 5 /dev/watchdog
}

