# cznic/watchdog_adjust/files/etc/init.d/watchdog_adjust
#!/bin/sh /etc/rc.common

START=99

start() {
	ubus call system watchdog '{ "frequency": '2' }'
	sleep 3	
	ubus call system watchdog '{ "timeout": '7' }'
}

