# collect/sentinel/sentinel-proxy/files/renew.cron
## crontab
#

# Periodically check and renew Sentinel:Proxy certificate
42 */12 * * * root sentinel-certgen certs --hooks-dir /usr/libexec/sentinel/renew_hooks.d
