# utils/logrotate/files/logrotate.cron
MAILTO=""
12     *       *       *       *       root	/usr/sbin/logrotate -s /tmp/logrotate.state /etc/logrotate.conf
