# net/squid/files/squid.cron
MAILTO=""
0	0	*	*	*	root	/usr/sbin/squid -k rotate
