# net/dnssec-rootkey/files/rootkeys_backup.cron
*	*	1	*	*	root	/usr/bin/rootkeys_backup.sh
