# hardware/atsha204/libatsha204/files/atsha204-feed-entropy.init
#!/bin/sh /etc/rc.common

START=01

start() {
	for i in $(seq 1 16); do
		/usr/bin/atsha204cmd feed-entropy
	done
}
