# cznic/updater/files/check-hashes.cron
5	0	*	*	*	root	touch /tmp/updater-check-hashes
