# cznic/cznic-cacert-bundle/files/get-crl.cron
MAILTO=""
*/30	*	*	*	*	root	/usr/bin/get-api-crl
