# cznic/userspace_time_sync/files/time_sync.init
#!/bin/sh /etc/rc.common

START=13

start() {
	if [ -c /dev/rtc0 ]; then
	        hwclock -u -s -f /dev/rtc0
	fi
}

stop() {
	if [ -c /dev/rtc0 ] ; then
		hwclock -u -w -f /dev/rtc0
	fi
}
