# admin/cronie/files/cron.init
#!/bin/sh /etc/rc.common

START=90

USE_PROCD=1

start_service() {
	mkdir -p /var/spool/cron
	ln -sf /etc/crontabs /var/spool/cron/crontabs

	procd_open_instance
	procd_set_param command /usr/sbin/crond -n
	procd_set_param respawn
	procd_close_instance
}
