# cznic/turris-netboot/files/netboot-update.cron
MAILTO=""

15	*/12	*	*	*	root	/usr/bin/netboot-manager update_rootfs -s
