# cznic/libatsha204/files/atsha204-feed-entropy.init
#!/bin/sh /etc/rc.common

START=11

start() {
	if [ -x /usr/bin/atsha204cmd ]; then
		for i in $(seq 1 16); do
			/usr/bin/atsha204cmd feed-entropy
		done
	fi
}

stop() {
	:
}
