# cznic/nuci/files/autocollect.cron
MAILTO=""
@daily	root	/usr/bin/nuci-helper-autocollect
