# cznic/ltemetr-core/files/ltemetr.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2019 CZ.NIC z.s.p.o. (https://www.nic.cz/)

START=99
USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command ltemetr
	procd_set_param respawn ${respawn_threshold:-3600} ${respawn_timeout:-5} ${respawn_retry:-5}
	procd_close_instance
}
