# cznic/netmetr/files/netmetr.cron
MAILTO=""
00 	* 	* 	* 	* 	root	netmetr --autostart --dwlhist --rwait 1800 &
