# net/nextcloud/files/nextcloud.cron
*/5  *  *  *  *  nobody /usr/bin/nextcloud-cron.sh
